# DEBIAN/postinst
#!/bin/bash

set -e

chmod a+x /etc/init.d/mtversion
chmod ug+rw /etc/mtversion
update-rc.d mtversion defaults
